# kanata.cron
# Add to root crontab with: sudo crontab -e
* * * * * pgrep -x kanata > /dev/null || sudo tmux new -d -s kanata '/opt/homebrew/bin/kanata --cfg /Users/tobyc/Documents/repos/scripts/src/kanata/glove80.kbd'
