# scripts/kanata.cron
# Add to root crontab with: sudo crontab -e
* * * * * pgrep -x kanata > /dev/null || /opt/homebrew/bin/kanata --cfg /Users/tcourtis/repos/keyboard-configs/glove80.kbd &
